# debian/axon-test/DEBIAN/postinst
sudo systemctl daemon-reload
sudo systemctl enable test-axon.service
sudo systemctl status test-axon.service
